# system/systemd/prasanna-theme.path
[Unit]
Description=Watch for Material You color changes

[Path]
PathModified=%h/.local/share/color-schemes/MaterialYouDark.colors

[Install]
WantedBy=default.target
